# config/rc.local
exec 2> /tmp/rc.local.log      # send stderr from rc.local to a log file
exec 1>&2                      # send stdout to the same log file
set -x                         # tell sh to display commands before execution

. /home/user/.virtualenvs/qwirk-python/bin/activate

cd /home/user/path/to/server

uwsgi --ini config/server_uwsgi.ini &

(python manage.py runworker &) && (daphne server.asgi:channel_layer &)

deactivate

exit 0
